# systemd/duskwatch-tray.service
[Unit]
Description=Duskwatch standalone tray app
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/Projects/duskwatch/tray-app/duskwatch_tray.py
Restart=on-failure

[Install]
WantedBy=graphical-session.target
